# manifest/etc/supervisor/01-php-fpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm7 -F
autostart=true
autorestart=true
stopsignal=INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=0
